# deploy/user-profile-analyzer.service
[Unit]
Description=User Profile Analyzer Web UI
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/user-profile-analyzer
Environment="PYTHONPATH=/opt/user-profile-analyzer/src"
Environment="APP_ENV=prod"
ExecStart=/opt/user-profile-analyzer/.venv/bin/python -m user_profile_analyzer.web_ui
Restart=always
RestartSec=5

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
